# supervisord.conf
[supervisord]
nodaemon=true

[program:webserver]
command=airflow webserver -p 8080
autostart=true
autorestart=true

[program:scheduler]
command=airflow scheduler
autostart=true
autorestart=true